# servicios/log-cleaner/systemd/robiotec-log-cleaner.service
[Unit]
Description=Robiotec weekly service log cleaner

[Service]
Type=oneshot
ExecStart=/root/robiotec/servicios/log-cleaner/scripts/clean-logs.sh
